# jade-desktop-mode.service
[Unit]
Description=Run JADE in Desktop Mode

[Service]
ExecStart=/usr/bin/jade-desktop-mode.sh

[Install]
WantedBy=multi-user.target 
